# deploy/80.ESP8266-Sensor.conf
[program:esp8266-uwsgi]
process_name=esp8266-uwsgi
command=/home/debian/.virtualenvs/ESP8266-Sensor/bin/uwsgi
    --chdir /web/ESP8266-Sensor/server
    --wsgi-file esp_server.py
    --callable app --processes 1 --threads 1
    --socket 127.0.0.1:3031 
environment=PATH="/home/myuser/.virtualenvs/myproj/bin:%(ENV_PATH)s"
